# packaging/watchguard.service

[Unit]
Description=Watchguard Host Health Monitor
Documentation=man:watchguard(8)
After=network-online.target
Wants=network-online.target

StartLimitIntervalSec=60
StartLimitBurst=10

[Service]
Type=simple

WorkingDirectory=/

ExecStart=/usr/bin/watchguard daemon --config /etc/watchguard/config.toml

Restart=on-failure
RestartSec=3

TimeoutStartSec=30
TimeoutStopSec=15

NoNewPrivileges=true

StandardOutput=journal
StandardError=journal

User=root
Group=root

[Install]
WantedBy=multi-user.target
